# filebeat/supervisord.conf
; Copyright (c) 2023 Battelle Energy Alliance, LLC.  All rights reserved.

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
environment=PATH=/usr/share/filebeat:%(ENV_PATH)s

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:filebeat]
command=/usr/local/bin/docker-entrypoint -e --strict.perms=false
user=%(ENV_PUSER)s
startsecs=0
startretries=0
stopasgroup=true
killasgroup=true
directory=/usr/share/filebeat
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:filebeat-nginx]
command=bash -c "/usr/local/bin/opensearch_status.sh -t malcolm_beats_template && /usr/local/bin/docker-entrypoint -e --strict.perms=false \
  --path.home /usr/share/filebeat-nginx \
  --path.config /usr/share/filebeat-nginx \
  --path.data /usr/share/filebeat-nginx/data \
  -c /usr/share/filebeat-nginx/filebeat-nginx.yml \
  --modules nginx"
user=%(ENV_PUSER)s
autostart=%(ENV_NGINX_LOG_ACCESS_AND_ERRORS)s
startsecs=30
startretries=2000000000
stopasgroup=true
killasgroup=true
directory=/usr/share/filebeat-nginx
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:filebeat-tcp]
command=bash -c "/usr/local/bin/opensearch_status.sh -t malcolm_beats_template && /usr/local/bin/docker-entrypoint -e --strict.perms=false \
  --path.home /usr/share/filebeat-tcp \
  --path.config /usr/share/filebeat-tcp \
  --path.data /usr/share/filebeat-tcp/data \
  -c /usr/share/filebeat-tcp/filebeat-tcp.yml"
user=%(ENV_PUSER)s
autostart=%(ENV_FILEBEAT_TCP_LISTEN)s
startsecs=30
startretries=2000000000
stopasgroup=true
killasgroup=true
directory=/usr/share/filebeat-tcp
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:watch-upload]
command=python3 /usr/local/bin/filebeat-watch-zeeklogs-uploads-folder.py %(ENV_PCAP_PIPELINE_VERBOSITY)s
  --start-sleep 30
  --polling "%(ENV_FILEBEAT_WATCHER_POLLING)s"
  --closed-sec %(ENV_FILEBEAT_WATCHER_POLLING_ASSUME_CLOSED_SEC)s
  --in "%(ENV_FILEBEAT_ZEEK_DIR)s"/upload
  --out "%(ENV_FILEBEAT_ZEEK_DIR)s"
  --uid %(ENV_PUID)s
  --gid %(ENV_PGID)s
user=root
startsecs=35
startretries=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:cron]
autorestart=true
command=/usr/local/bin/supercronic -json "%(ENV_SUPERCRONIC_CRONTAB)s"
user=%(ENV_PUSER)s
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
