# netbox/supervisord.conf
; Copyright (c) 2022 Battelle Energy Alliance, LLC.  All rights reserved.

[inet_http_server]
port=0.0.0.0:9001

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001

[group:netbox]
programs=main,worker,housekeeping,initialization

[program:main]
command=/opt/netbox/launch-netbox.sh
autostart=true
autorestart=true
startsecs=180
stopwaitsecs=30
startretries=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
redirect_stderr=true

[program:initialization]
command=/opt/netbox/venv/bin/python /usr/local/bin/netbox_init.py
  --wait
  --url "http://localhost:8080/netbox"
  --token "%(ENV_SUPERUSER_API_TOKEN)s"
  --net-map /usr/local/share/net-map.json
autostart=true
autorestart=false
startsecs=0
startretries=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
redirect_stderr=true

[program:worker]
command=/opt/netbox/venv/bin/python /opt/netbox/netbox/manage.py rqworker
autostart=true
autorestart=true
startsecs=60
stopwaitsecs=30
startretries=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
redirect_stderr=true

[program:housekeeping]
command=/opt/netbox/housekeeping.sh
autostart=true
autorestart=true
startsecs=60
stopwaitsecs=30
startretries=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
redirect_stderr=true

[program:cron]
command=/usr/local/bin/supercronic -json "%(ENV_SUPERCRONIC_CRONTAB)s"
autostart=%(ENV_NETBOX_CRON)s
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
