# name-map-ui/config/supervisord.conf
; Copyright (c) 2022 Battelle Energy Alliance, LLC.  All rights reserved.

[unix_http_server]
file=/tmp/supervisor-main.sock   ; (the path to the socket file)
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisor-main.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor-main.sock

[program:php-fpm]
command=/usr/sbin/php-fpm8 -F
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=false
startretries=0

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=false
startretries=0

[program:logaccess]
command=/usr/bin/tail -F /var/log/nginx/access.log
startsecs=10
startretries=2000000000
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stdout_logfile_maxbytes=0
redirect_stderr=false

[program:logerrors]
command=/usr/bin/tail -F /var/log/nginx/error.log
startsecs=10
startretries=2000000000
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/null
stdout_logfile_maxbytes=0
redirect_stderr=false

[program:watch-upload]
command=/bin/bash -c "sleep 15 && /usr/local/bin/name-map-save-watch.sh"
startsecs=20
startretries=1
stopasgroup=true
killasgroup=true
directory=/var/www/html/upload
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
