# pcap-monitor/supervisord.conf
; Copyright (c) 2022 Battelle Energy Alliance, LLC.  All rights reserved.

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:watch-upload]
command=/bin/bash -c "sleep 30 && /usr/local/bin/watch-pcap-uploads-folder.sh"
startsecs=35
startretries=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:pcap-publisher]
command=python3 /usr/local/bin/pcap_watcher.py
  --verbose "%(ENV_PCAP_PIPELINE_DEBUG)s"
  --extra-verbose "%(ENV_PCAP_PIPELINE_DEBUG_EXTRA)s"
  --opensearch "%(ENV_OPENSEARCH_URL)s"
  --opensearch-wait
  --node "%(ENV_PCAP_NODE_NAME)s"
  --ignore-existing "%(ENV_PCAP_PIPELINE_IGNORE_PREEXISTING)s"
  --start-sleep 60
  --directory "%(ENV_PCAP_PATH)s"/processed
user=%(ENV_PUSER)s
startsecs=65
startretries=1
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
