# sensor-raspi/rootfs/etc/systemd/system/rpi-set-sysconf.service
[Unit]
Description=Set up system configuration

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/rpi-set-sysconf

[Install]
RequiredBy=basic.target
