# logstash/supervisord.conf
; Copyright (c) 2025 Battelle Energy Alliance, LLC.  All rights reserved.

[inet_http_server]
port=0.0.0.0:9001

[supervisord]
nodaemon=true
loglevel=error
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
environment=PATH=/usr/share/logstash/bin:%(ENV_PATH)s

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001

[program:logstash]
command=/usr/local/bin/logstash-start.sh
autostart=true
autorestart=true
startsecs=0
startretries=0
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stopsignal=INT
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
