# strelka/backend/supervisord.conf
; Copyright (c) 2025 Battelle Energy Alliance, LLC.  All rights reserved.

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
chmod=0700

[supervisord]
nodaemon=true
loglevel=error
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:backend]
command=/home/strelka/.pyenv/bin/pyenv exec strelka-backend
numprocs=%(ENV_STRELKA_BACKEND_PROCS)s
process_name=%(program_name)s_%(process_num)d
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:cron]
command=/usr/local/bin/supercronic -json "%(ENV_SUPERCRONIC_CRONTAB)s"
numprocs=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:ruleinit]
command=/usr/local/bin/yara_rules_setup.sh -s
numprocs=1
autostart=true
autorestart=false
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:freshclam]
command=/usr/bin/freshclam freshclam --config-file=/etc/clamav/freshclam.conf --daemon
numprocs=1
autostart=%(ENV_RULES_UPDATE_ENABLED)s
autorestart=%(ENV_RULES_UPDATE_ENABLED)s
startsecs=0
startretries=0
stopasgroup=true
killasgroup=true
directory=/
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:clamd]
command=/usr/sbin/clamd -c /etc/clamav/clamd.conf
numprocs=1
autostart=true
autorestart=true
startsecs=0
startretries=0
stopasgroup=true
killasgroup=true
directory=/
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true