# pcap-capture/supervisord.conf
; Copyright (c) 2025 Battelle Energy Alliance, LLC.  All rights reserved.

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
chmod=0700

[supervisord]
nodaemon=true
user=root
loglevel=error
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:netsniff-roll]
startsecs=15
command=/usr/local/bin/netsniff-roll.sh
user=%(ENV_PUSER)s
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=%(ENV_PCAP_ENABLE_NETSNIFF)s
autorestart=%(ENV_PCAP_ENABLE_NETSNIFF)s
directory=%(ENV_PCAP_PATH)s

[include]
files = /etc/supervisor.d/*.conf
