# pcap-capture/supervisord.conf
; Copyright (c) 2019 Battelle Energy Alliance, LLC.  All rights reserved.

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:netsniff-roll]
startsecs=15
command=/usr/local/bin/netsniff-roll.sh
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=%(ENV_PCAP_ENABLE_NETSNIFF)s
directory=%(ENV_PCAP_PATH)s
user=%(ENV_PCAP_USER)s

[include]
files = /etc/supervisor.d/*.conf
